# files/etc/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:boinc]
#command=/usr/bin/b i
command=boinc --allow_remote_gui_rpc --dir $BOINC_DATADIR
autorestart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:gridcoinresearch]
#command=/usr/bin/grc i
command=/usr/bin/gridcoinresearchd -server -shrinkdebugfile -upnp -blocknotify='exec /grcupdate.sh'
autorestart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
